# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:lighttpd]
command=/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
redirect_stderr=true
environment=HOME="/tmp"

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/uwsgi.adagios.ini --fastcgi-socket /tmp/adagios.socket -L
redirect_stderr=true
environment=HOME="/tmp"

[program:nagios]
command=/usr/sbin/nagios /etc/nagios/nagios.cfg

[program:npcd]
command=/usr/sbin/npcd -f /etc/pnp4nagios/npcd.cfg

[group:nagios]
programs=uwsgi,lighttpd,nagios,npcd

[program:postfix]
process_name = master
directory = /etc/postfix
command = /usr/bin/pidproxy /var/spool/postfix/pid/master.pid /usr/sbin/postfix -c /etc/postfix start
startsecs = 0
autorestart = false

[program:nagiospid]
command=/usr/bin/nagios-supervisor-wrapper.sh genpid
autorestart=false
